# packages/asus-charge-daemon/asus-charge-daemon-0.2.0.tar.gz/asus-charge-daemon-0.2.0/data/asuscharged.service
[Unit]
Description=ASUS Battery Charge Control Daemon

[Service]
Type=dbus
BusName=ca.cforrester.AsusChargeDaemon1
Restart=always
RestartSec=10
User=root
StateDirectory=asuscharged
ConfigurationDirectory=asuscharged
Environment="PYTHONOPTIMIZE=1"
ExecStart=/usr/bin/asuscharged

[Install]
WantedBy=multi-user.target
