# packages/nexpose-py/nexpose-py-0.0.18.tar.gz/nexpose-py-0.0.18/extras/nexpose-exporter.service
[Unit]
Description=Monitor Nexpose scans for Prometheus export
Documentation=https://github.com/coalfire/nexpose-py
After=syslog.target
After=network.target

[Service]
Type=simple
User=prometheus
Group=prometheus
EnvironmentFile=/etc/default/nexpose-exporter
ExecStart=/usr/local/bin/nsc-exporter $ARGS
Restart=always
RestartSec=42sec

[Install]
WantedBy=multi-user.target
